# udev-rules/91-power-manager.rules.in
#input(like bt)
ACTION=="add" SUBSYSTEM=="input"  DEVPATH=="*/input[0-9]*/event[0-9]*" RUN+="@PREFIX@/bin/pm_event add $DEVNAME"
ACTION=="remove" SUBSYSTEM=="input"  DEVPATH=="*/input[0-9]*/event[0-9]*" RUN+="@PREFIX@/bin/pm_event remove $DEVNAME"

# Since power_manager is running as app user, make sure this user can read the input device nodes
SUBSYSTEM=="input" DEVPATH=="*/input[0-9]*/event[0-9]*" GROUP="video"

