# udev-rules/91-power-manager.rules.in
#input(like bt)
ACTION=="add" SUBSYSTEM=="input"  DEVPATH=="*/input[1-9]*/event[1-9]*" RUN+="@PREFIX@/bin/pm_event add $DEVNAME"
ACTION=="remove" SUBSYSTEM=="input"  DEVPATH=="*/input[1-9]*/event[1-9]*" RUN+="@PREFIX@/bin/pm_event remove $DEVNAME"
